# systemd/frankfurt-bot.service
[Unit]
Description=Frankfurt discord bot

[Service]
WorkingDirectory=/home/njones/frankfurt/
ExecStart=/home/njones/frankfurt/run-bot.sh

[Install]
WantedBy=default.target
